# debian/postinst
#!/bin/bash

set -e

update-rc.d cloudshell2-fan defaults
